# ansible/roles/wake_on_lan/templates/wol-enable.service.j2
[Unit]
Description=Enable Wake-up on LAN for {{ item }}

[Service]
Type=oneshot
ExecStart=/sbin/ethtool -s {{ item }} wol g

[Install]
WantedBy=basic.target
